# scripts/pnrs.service
[Unit]
Description=PNRS - Private NPM Registry Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=pnrs
Group=pnrs
WorkingDirectory=/opt/pnrs
ExecStart=/opt/pnrs/pnrs
Restart=always
RestartSec=5

# Environment variables
Environment=PNRS_HOST=0.0.0.0
Environment=PNRS_PORT=8000
Environment=PNRS_UPSTREAM_REGISTRY=https://registry.npmjs.org
Environment=RUST_LOG=info

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/pnrs/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
